# filesplitter.service
[Unit]
Description=File Splitter Flask Application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/app
Environment="PATH=/path/to/your/venv/bin"
Environment="PYTHONPATH=/path/to/your/app"
EnvironmentFile=/path/to/your/app/.env
ExecStart=/path/to/your/venv/bin/python run.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target